# Openwrt/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
export MQTTLAUNCHERCONFIG=/root/launcher.conf
/root/mqtt-launcher.py > /var/log/mqtt-launcher 2>&1
exit 0
